# relay/deploy/cockpit-relay.service
[Unit]
Description=Cockpit Relay Server
After=network.target

[Service]
Type=simple
User=len
WorkingDirectory=/home/len/projects/claude-cockpit/relay
EnvironmentFile=/home/len/projects/claude-cockpit/relay/.env
ExecStart=/home/len/projects/claude-cockpit/relay/venv/bin/python -m relay.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
